# dynamic_proxy/supervisord.conf
[supervisord]
nodaemon=true

[program:dynaproxy]
command=/usr/bin/python /root/dynaproxy.py
directory=/root
stdout_logfile=/var/log/dynaproxy.log
redirect_stderr=true
autorestart=true
stopsignal=INT
